# src/bios-agent-smtp.service.in
[Unit]
Description=bios-agent-smtp service
Requires=network.target malamute.service
After=network.target malamute.service
Before=bios-agent-asset.service
PartOf=bios.target

[Service]
Type=simple
User=bios
EnvironmentFile=-/usr/share/bios/etc/default/bios
EnvironmentFile=-/usr/share/bios/etc/default/bios__%n.conf
EnvironmentFile=-/etc/default/bios
EnvironmentFile=-/etc/default/bios__%n.conf
ExecStart=@prefix@/bin/bios-agent-smtp
Restart=always

[Install]
WantedBy=bios.target
