# src/bios-agent-smtp.service.in
[Unit]
Description=bios-agent-smtp service
Requires=network.target malamute.service
After=network.target malamute.service
Before=bios-agent-asset.service
PartOf=bios.target

[Service]
Type=simple
User=bios
ExecStart=@prefix@/bin/bios-agent-smtp --config /etc/agent-smtp/bios-agent-smtp.cfg
Restart=always

[Install]
WantedBy=bios.target
