# src/fty-email.service.in
[Unit]
Description=fty-email service
Requires=network.target malamute.service
After=network.target malamute.service
Before=fty-asset.service
PartOf=bios.target

[Service]
Type=simple
User=bios
EnvironmentFile=-/usr/share/bios/etc/default/bios
EnvironmentFile=-/etc/default/bios
Environment="prefix=@prefix@"
ExecStart=@prefix@/bin/fty-email --config @sysconfdir@/@PACKAGE@/fty-email.cfg
Restart=always

[Install]
WantedBy=bios.target
